# deploy/solarica-backend.service
[Unit]
Description=Solarica FastAPI Backend
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=exec
User=solarica
Group=solarica
WorkingDirectory=/opt/solarica/backend
EnvironmentFile=/opt/solarica/.env
ExecStart=/opt/solarica/venv/bin/uvicorn app.main:app \
    --host 127.0.0.1 \
    --port 8000 \
    --workers 2 \
    --no-access-log

# Restart policy
Restart=always
RestartSec=5
StartLimitInterval=60
StartLimitBurst=3

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/opt/solarica

[Install]
WantedBy=multi-user.target
