# Solarica/scripts/systemd/solarica.service
# /etc/systemd/system/solarica.service
#
# Solarica FastAPI backend — native systemd unit.
# User: solarica   Port: 8013   App path: /home/solarica/app/backend
#
# Install:
#   sudo cp solarica.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now solarica
#
# Log tail:   journalctl -u solarica -f
# Status:     systemctl status solarica
# Health:     curl -s http://127.0.0.1:8013/api/health

[Unit]
Description=Solarica FastAPI Backend (solarica.wavelync.com)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=solarica
Group=solarica
WorkingDirectory=/home/solarica/app/backend

# Load app-level environment variables (DATABASE_URL, SECRET_KEY, etc.)
EnvironmentFile=/home/solarica/app/backend/.env
Environment=PYTHONPATH=/home/solarica/app/backend

ExecStart=/home/solarica/app/venv/bin/uvicorn main:app \
    --host 127.0.0.1 \
    --port 8013 \
    --workers 1

Restart=always
RestartSec=5

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
